# srcs/nginx/supervisor.ini
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=sh -c "nginx -g 'daemon off;'"
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true

[program:ssh]
command=sh -c "/usr/sbin/sshd -D"
autostart=true
autorestart=true
